# systemd/hey-claude.service.in
[Unit]
Description=hey-claude — push-to-talk Claude voice assistant daemon

[Service]
Type=dbus
BusName=org.kdc.HeyClaude
ExecStart=@HEYCLAUDED@
Restart=on-failure

[Install]
WantedBy=default.target
